# deployment/supervisor/paradisedollz-reverb.conf
[program:paradisedollz-reverb]
; Keeps the Laravel Reverb WebSocket server alive on 127.0.0.1:8080.
; Nginx proxies wss://yourdomain.com/app/* to this process.
; If it crashes, Supervisor restarts it automatically.

process_name=%(program_name)s
command=php /var/www/paradisedollz/artisan reverb:start --host=127.0.0.1 --port=8080 --no-ansi
autostart=true
autorestart=true
stopwaitsecs=10
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/paradisedollz-reverb.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
