# deployment/supervisor/paradisedollz-queue.conf
[program:paradisedollz-queue]
; Processes queued jobs from the database jobs table:
;   - All 7 email types (Mailables with ShouldQueue)
;   - Async broadcast events (ShouldBroadcast) for community chat
;
; --sleep=3   : poll interval when queue is empty (seconds)
; --tries=3   : retry a failed job up to 3 times before logging to failed_jobs
; --timeout=60: kill a job that runs longer than 60 seconds

process_name=%(program_name)s_%(process_num)02d
command=php /var/www/paradisedollz/artisan queue:work database --sleep=3 --tries=3 --timeout=60 --no-ansi
autostart=true
autorestart=true
stopwaitsecs=30
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/paradisedollz-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
